# scheduler-worktree/cleanup_assets.timer
[Unit]
Description=Periodic assets duplicate cleanup
Documentation=Cleans up duplicate images in assets folder that already exist in cache
After=docker.service

[Service]
Type=oneshot
User=warby
WorkingDirectory=/home/warby/Workspace/jupyter
ExecStart=/usr/bin/docker compose exec -T jupyterlab python /home/jovyan/workspaces/cleanup_assets_duplicates.py --cleanup --quiet
StandardOutput=journal
StandardError=journal

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
